# src/etc/s6-overlay/s6-rc.d/generate-ssh-keys/up
/etc/s6-overlay/scripts/generate-ssh-keys